# deploy/kanalen/launchagents/dev.activechain.kanalen.round.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>dev.activechain.kanalen.round</string>
  <key>ProgramArguments</key><array>
    <string>/Users/johansellstrom/activechain-deploy/kanalen/current/scripts/run-kanalen-round.sh</string>
  </array>
  <key>RunAtLoad</key><true/><key>StartInterval</key><integer>30</integer>
  <key>StandardOutPath</key><string>/Users/johansellstrom/activechain-deploy/kanalen/logs/round.log</string>
  <key>StandardErrorPath</key><string>/Users/johansellstrom/activechain-deploy/kanalen/logs/round.log</string>
</dict></plist>
