# deploy/kanalen/launchagents/dev.activechain.kanalen.rpc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.activechain.kanalen.rpc</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Volumes/ActiveChain/testnet/current/bin/activechain-rpc-node</string>
    <string>/Volumes/ActiveChain/testnet/rpc/rpc-index.snapshot</string>
    <string>0.0.0.0:49151</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/johansellstrom</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/johansellstrom/Library/Logs/ActiveChain/rpc.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/johansellstrom/Library/Logs/ActiveChain/rpc.err.log</string>
</dict>
</plist>
